# i386/Packages/Extra/wifibox_1.0.4_i386/DEBIAN/postinst
#!/bin/sh
# postinst script for wifi-box
#
# see: dh_installdeb(1)

set -e

. /usr/share/debconf/confmodule

if [ -f /etc/network.conf.bak ]; then
	mv -f /etc/network.conf.bak /etc/network.conf
	[ -f /etc/wpa/wpa.conf.bak ] && mv -f /etc/wpa/wpa.conf.bak /etc/wpa/wpa.conf
fi

[ -d /etc/frisbee ] && rm -fr /etc/frisbee
[ -f /etc/systemd/system/frisbee.service ] && rm -f /etc/systemd/system/frisbee.service
# dhcpcd -b -d

if [ -x "/etc/init.d/rcS_slitaz" ]; then
	update-rc.d rcS_slitaz defaults >/dev/null
	if [ -n "$2" ]; then
		_dh_action=restart
	else
		_dh_action=start
	fi
	invoke-rc.d rcS_slitaz $_dh_action || exit $?
fi

if [ -n "`pidof network_tray`" ]; then
killall network_tray
fi

[ -f /tmp/instyad.sh ] && nohup /tmp/instyad.sh > /tmp/nohup-wbox.out 2>&1&
if [ -z "`pidof netmon_wce`" ]; then
nohup netmon_wce > /tmp/nohup-wifibox.out 2>&1&
fi

db_stop;

exit 0
